# infrastructure/ansible/roles/patroni/templates/patroni.service.j2
[Unit]
Description=Patroni PostgreSQL cluster manager
After=network.target

[Service]
Type=simple
User=postgres
ExecStart=/opt/patroni-venv/bin/patroni /etc/patroni.yml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target